# yuzu/yuzu-update.timer
[Unit]
Description=Update yuzu every hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target